# roles/inbox/templates/document-scanner.service.j2
[Unit]
Description=Automatic document scanner

[Service]
WorkingDirectory={{ inbox_root }}/scanadf
ExecStart={{ inbox_bindir }}/document-scanner

[Install]
WantedBy=basic.target
